# deploy/services/nginx/run
#!/bin/bash -e
exec log-service nginx -- nginx -c /etc/nginx/nginx.conf
