# deploy/services/icecast/run
#!/bin/bash -e
exec log-service icecast -- icecast2 -c /etc/icecast2/icecast.xml
